# build/iso_root/etc/inittab
# /etc/inittab
id:3:initdefault:
ca::ctrlaltdel:/sbin/shutdown -t3 -r now
