# services/system/carrot-shell.service
[Unit]
Description=CarrotOS Desktop Environment
After=display-manager.service
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/carrot-shell
Restart=on-failure
RestartSec=3

# Environment
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_SESSION_TYPE=wayland
Environment=QT_QPA_PLATFORM=wayland

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Security
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=graphical-session.target
